# scripts/apply-smart-backend.service
[Unit]
Description=Apply Smart AI Backend
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/apply-smart-ai

# Environment variables (load from .env file)
EnvironmentFile=/var/www/apply-smart-ai/.env

# Java options for production
Environment="JAVA_OPTS=-Xmx1024m -Xms512m -XX:+UseG1GC -XX:MaxGCPauseMillis=200"

# Spring Boot profile
Environment="SPRING_PROFILES_ACTIVE=prod"

# Execute the application
ExecStart=/usr/bin/java $JAVA_OPTS -jar /var/www/apply-smart-ai/app.jar

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=apply-smart-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
